# boards/iconia-w4-820/install/iconia-powertune.conf
# iconia-powertune.conf — apply UX-safe power savings at every boot.
# Pins schedutil, disables the broken idle-suspend (screen-off still allowed),
# and enables light I2C runtime-PM. See iconia-powertune.sh for rationale.
description "Iconia W4-820 UX-safe power tuning"
author "iconia"

start on started system-services
task

exec /usr/local/sbin/iconia-powertune
