# boards/iconia-w4-820/install/iconia-install.conf
# Iconia W4-820 — one-shot eMMC auto-install trigger.
# Fires only when the boot cmdline contains `iconia_install=1`, after
# system-services are up (chromeos-install needs udev/cros-disks/initctl).
#
# NOT a `task`: we launch the installer DETACHED (setsid ... &) and return
# immediately, so the job does NOT block the boot milestone chain (a blocking
# task here wedged the boot at the splash). Boot proceeds to OOBE while the
# install runs underneath and powers off when done.
description "Iconia one-shot eMMC auto-install"
author "iconia project"

start on started system-services

script
  if grep -qw iconia_install /proc/cmdline; then
    setsid /sbin/iconia-install.sh </dev/null >/dev/null 2>&1 &
  fi
end script
