# boards/iconia-w4-820/install/iconia-buttond.conf
# iconia-buttond.conf — long-press the Windows/home button to open crosh.
# Watches the soc_button_array evdev node and injects Ctrl+Alt+T via uinput on a
# long KEY_LEFTMETA hold. Short press keeps its normal "home" action.
description "Iconia Windows-button long-press -> crosh"
author "iconia"

start on started ui
stop on stopping ui
respawn
respawn limit 10 5

pre-start exec modprobe uinput
exec /usr/local/sbin/iconia-buttond
