# boards/iconia-w4-820/install/iconia-shutdown-black.conf
# iconia-shutdown-black — kill the lingering light at shutdown/reboot (retina-safe).
#
# PARTIAL fix. The bright "fade to white" at shutdown is ash's shutdown animation,
# which is COMPILED INTO CHROME and plays BEFORE any external hook (powerd, Upstart)
# can run — ChromeOS does animation first, THEN issues the shutdown request. So we
# CANNOT recolor or suppress that ~1s flash without rebuilding Chrome (not worth it).
# What this job DOES kill is the multi-second white that used to linger through the
# unmount/power-down AFTER the animation: it fires once Chrome has exited, stops
# powerd (so its fade can't re-raise the backlight), paints a solid black frame with
# frecon (Chrome released DRM master by then), and cuts the backlight. Proven to hold
# the panel black for the whole post-animation window.
#
# The residual animation flash is made non-blinding at night by the lowered ALS floor
# (5%) in iconia-als-brightness-install.sh — the two together are the accepted solution.
start on starting halt or starting reboot
task
script
  bl=/sys/class/backlight/intel_backlight
  # wait up to ~3s for Chrome to release the display (DRM master) and exit
  i=0; while pgrep -x chrome >/dev/null 2>&1 && [ $i -lt 15 ]; do sleep 0.2; i=$((i+1)); done
  stop powerd 2>/dev/null || true
  frecon --clear 0x000000 --daemon --no-login --num-vts 1 2>/dev/null || true
  # hammer the backlight off for ~1.6s to outlast any residual powerd fade
  j=0; while [ $j -lt 8 ]; do echo 4 > "$bl/bl_power" 2>/dev/null; echo 0 > "$bl/brightness" 2>/dev/null; sleep 0.2; j=$((j+1)); done
end script
