# boards/iconia-w4-820/install/iconia-buttond.conf
# iconia-buttond.conf — long-press the Windows/home button to open crosh.
# Watches the soc_button_array evdev node and injects Ctrl+Alt+T via uinput on a
# long KEY_LEFTMETA hold. Short press keeps its normal "home" action.
description "Iconia Windows-button long-press -> crosh"
author "iconia"

start on started ui
stop on stopping ui
respawn
respawn limit 10 30

# S20: on 6.6.99 soc_button_array is a *module*, so its two "gpio-keys" evdev
# nodes appear slightly AFTER `started ui`. iconia-buttond scans /dev/input once
# and exits if no KEY_LEFTMETA device exists yet — which burned the respawn
# budget at boot (service ended up stop/waiting, feature dead). Block in
# pre-start until the gpio-keys input device is registered before exec'ing.
pre-start script
  modprobe uinput || true
  for i in $(seq 1 60); do
    if grep -q 'Name="gpio-keys"' /proc/bus/input/devices 2>/dev/null; then
      exit 0
    fi
    sleep 1
  done
  logger -t iconia-buttond "gpio-keys evdev never appeared after 60s"
  exit 0
end script

exec /usr/local/sbin/iconia-buttond
