# boards/iconia-w4-820/install/iconia-install.conf
# Iconia W4-820 — one-shot eMMC auto-install trigger.
# Fires only when the boot cmdline contains `iconia_install=1`, after
# system-services are up (chromeos-install needs udev/cros-disks/initctl).
# The script itself is guarded by a sentinel on the USB ESP so it runs once.
description "Iconia one-shot eMMC auto-install"
author "iconia project"

start on started system-services
task

script
  if grep -qw iconia_install /proc/cmdline; then
    /sbin/iconia-install.sh
  fi
end script
