# boards/iconia-w4-820/install/iconia-memtune.conf
# iconia-memtune.conf — apply low-RAM memory tuning at boot (zram->zstd + vm knobs).
# Runs after system-services so ChromeOS swap setup is already up (we convert it).
description "Iconia W4-820 memory tuning (2GB)"
author "iconia"

start on started system-services
task

exec /usr/local/sbin/iconia-memtune
