# boards/iconia-w4-820/install/iconia-buttond.conf
# iconia-buttond.conf — remap the Windows/home button.
# Watches the soc_button_array evdev node (KEY_LEFTMETA) and adds two gestures:
#   long hold  -> inject Ctrl+Alt+T via uinput (open crosh)
#   double-tap -> powerd_dbus_suspend (the only on-demand sleep path; see the
#                 header of iconia-buttond.c — RTC wake is dead and Chrome's
#                 tablet power-menu has no "Sleep", but s2idle + power-button
#                 wake work). Short press keeps its normal "home" action.
description "Iconia Windows-button: long-press -> crosh, double-tap -> sleep"
author "iconia"

start on started ui
stop on stopping ui
respawn
respawn limit 10 30

# S20: two boot-time races have to clear before buttond can work, so we block in
# pre-start until BOTH are satisfied:
#   1) INPUT (read side): on 6.6.99 soc_button_array is a *module*, so its
#      "gpio-keys" evdev nodes register slightly after `started ui`. buttond
#      scans /dev/input once and exits if no KEY_LEFTMETA device exists yet.
#   2) OZONE (inject side): buttond creates a uinput keyboard; ChromeOS/ozone
#      only routes injected accelerators from a device it enumerated AFTER its
#      input stack was up. If the uinput device is created too early at boot,
#      ozone never adopts it and the injected Ctrl+Alt+T goes nowhere (the
#      service looked "running" but the long-press did nothing). Waiting for a
#      chrome renderer process is a reliable "ozone is watching udev now" signal.
pre-start script
  modprobe uinput || true
  # (1) wait for the gpio-keys input device
  for i in $(seq 1 60); do
    grep -q 'Name="gpio-keys"' /proc/bus/input/devices 2>/dev/null && break
    sleep 1
  done
  # (2) wait for Chrome's input stack (a renderer implies ozone is fully up),
  # then a short settle so our UI_DEV_CREATE lands as a hot-plug ozone adopts.
  for i in $(seq 1 60); do
    pgrep -f 'chrome --type=renderer' >/dev/null 2>&1 && break
    sleep 1
  done
  sleep 5
end script

exec /usr/local/sbin/iconia-buttond
