# boards/iconia-w4-820/install/iconia-binfmt-misc.conf
# iconia-binfmt-misc.conf — mount binfmt_misc so ARC++ can register houdini.
#
# Legacy ARC++ (container Android) needs the ARM native-bridge (libhoudini.so)
# registered via binfmt_misc so all of Android's ARM /system/bin can execute on
# this x86 Bay Trail device. The kernel has CONFIG_BINFMT_MISC=y, but nothing on
# this build mounts the filesystem, so arc-setup's RegisterAllBinFmtMiscEntries
# fails (No such file: /proc/sys/fs/binfmt_misc/arm_exe|arm64_exe) and
# /system/bin/arcbootcontinue can't exec -> Android tears down at boot.
# This job mounts it early so houdini can register on each ARC start.
# See PROGRESS.md Session 15 + memory iconia-android-arc-diag.
description "mount binfmt_misc for ARC++ houdini native-bridge"
author "iconia"

start on started boot-services
task

script
  grep -qw binfmt_misc /proc/mounts || \
    mount -t binfmt_misc binfmt_misc /proc/sys/fs/binfmt_misc
end script
