# boards/iconia-w4-820/install/iconia-install.conf
# Iconia W4-820 — one-shot eMMC auto-install trigger.
# Fires only when the boot cmdline contains `iconia_install=1`.
#
# Trigger EARLY (started boot-services), NOT system-services: on this board the
# later service layer (shill/btmanagerd) crash-loops because their kernel
# modules are missing, so `started system-services` may never settle and a job
# waiting on it never fires. boot-services completes well before that storm.
#
# Launched DETACHED (setsid ... &) so it never blocks the boot milestone chain.
description "Iconia one-shot eMMC auto-install"
author "iconia project"

start on started boot-services

script
  if grep -qw iconia_install /proc/cmdline; then
    setsid /sbin/iconia-install.sh </dev/null >/dev/null 2>&1 &
  fi
end script
