# boards/iconia-w4-820/install/iconia-acpoll.conf
# iconia-acpoll.conf — keep AC/battery UI responsive after the GPE-storm fix.
# The kernel cmdline mutes the storming EC interrupt (gpiolib_acpi.ignore_interrupt=
# INT33FC:02@18), which removes the power_supply push-notify. This respawning job
# pokes powerd via a udev 'change' every few seconds so the tray tracks plug/unplug.
# See iconia-acpoll.sh for the full rationale.
description "Iconia W4-820 power_supply UI poke"
author "iconia"

start on started system-services
stop on stopping system-services
respawn

exec /usr/local/sbin/iconia-acpoll 5
